# dopamine.service
[Unit]
Description=Dopamine Hardware Bridge
After=network.target

[Service]
User=pi
WorkingDirectory=/home/pi/dopamine-hardware
Environment="PATH=/home/pi/.local/bin:/usr/bin"
ExecStart=/usr/bin/python3 app.py
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
